# worker/worker.target
[Unit]
Description=Worker Target
RefuseManualStart=yes
Wants=buildbot-worker.service
